# build/Info.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>CFBundleExecutable</key> <string>Mongoose</string>
  <key>CFBundlePackageType</key> <string>APPL</string>
  <key>CFBundleTypeRole</key> <string>None</string>
  <key>CFBundleIconFile</key> <string>mongoose.icns</string>
  <key>LSUIElement</key> <true/>
  <key>RunAtLoad</key> <true/>
  <key>Label</key> <string>com.kolkin.mongoose</string>
  <key>ProgramArguments</key> <array> </array>
  <key>KeepAlive</key> <true/>
</dict>
</plist>
